# auto-temp-ctrl/DEBIAN/postinst
#!/bin/bash
systemctl daemon-reload
systemctl enable auto-temp-ctrl.service
systemctl start auto-temp-ctrl.service
