# auto-temp-ctrl/DEBIAN/postinst
#!/bin/bash
set -e

VENV_DIR="/opt/auto-temp-ctrl-venv"
PY_MODULE="auto_temp_gpiod_ctrl"

# Check libgpiod
REQUIRED_GPIOD_VER="2.1.3"

check_libgpiod() {
    if ! command -v gpioinfo >/dev/null 2>&1; then
        echo "❌ 'gpioinfo' not found. libgpiod >= $REQUIRED_GPIOD_VER is required."
        echo "Please install libgpiod >= $REQUIRED_GPIOD_VER before installing auto-temp-ctrl."
        echo "For systems without a package, build it from source:"
        echo "  cd ~"
        echo "  wget https://mirrors.edge.kernel.org/pub/software/libs/libgpiod/libgpiod-2.1.3.tar.xz"
        echo "  tar -xf libgpiod-2.1.3.tar.xz"
        echo "  cd libgpiod-2.1.3"
        echo "  ./configure --enable-tools"
        echo "  make"
        echo "  sudo make install"
        echo "  sudo ldconfig"
        exit 1
    fi

    ver_raw=$(gpioinfo --version 2>&1 | head -n1)
    ver=$(echo "$ver_raw" | grep -Eo 'v[0-9]+\.[0-9]+\.[0-9]+' | head -n1 | sed 's/^v//')
    echo "Detected gpioinfo version: $ver_raw"
    if [[ -z "$ver" ]]; then
        echo "❌ Unable to determine gpioinfo version."
        exit 1
    fi
    if dpkg --compare-versions "$ver" lt "$REQUIRED_GPIOD_VER"; then
        echo "❌ Detected gpioinfo version $ver. Version >= $REQUIRED_GPIOD_VER is required."
        echo "Please upgrade libgpiod to >= $REQUIRED_GPIOD_VER before installing auto-temp-ctrl."
        exit 1
    else
        echo "✅ gpioinfo version $ver is sufficient (>= $REQUIRED_GPIOD_VER)"
    fi
}

check_libgpiod

# Setup venv and pip install
if [ ! -d "$VENV_DIR" ]; then
    echo "Creating Python venv at $VENV_DIR..."
    python3 -m venv "$VENV_DIR"
    echo "✅ Created Python venv at $VENV_DIR."
else
    echo "✅ Python venv already exists at $VENV_DIR."
fi

source "$VENV_DIR/bin/activate"

pip install --upgrade pip && echo "✅ pip is up to date."
pip install --upgrade --force-reinstall "$PY_MODULE" && echo "✅ $PY_MODULE is installed in the venv."

# Config file
if [ ! -f /etc/auto-temp-ctrl/auto-temp-ctrl.conf ]; then
    cp /usr/share/auto-temp-ctrl/auto-temp-ctrl.config /etc/auto-temp-ctrl/auto-temp-ctrl.conf
    echo "✅ Copied default config to /etc/auto-temp-ctrl/auto-temp-ctrl.conf."
else
    echo "✅ Config already exists at /etc/auto-temp-ctrl/auto-temp-ctrl.conf."
fi

# Systemd service
echo "Installing systemd service..."
cp /usr/share/auto-temp-ctrl/auto-temp-ctrl.service /etc/systemd/system/
if [ -f /etc/systemd/system/auto-temp-ctrl.service ]; then
    echo "✅ Service file installed to /etc/systemd/system/auto-temp-ctrl.service."
else
    echo "❌ Failed to copy auto-temp-ctrl.service!"
    exit 1
fi

systemctl daemon-reload && echo "✅ systemd daemon reloaded."
systemctl enable auto-temp-ctrl.service && echo "✅ Service enabled."
systemctl restart auto-temp-ctrl.service && echo "✅ Service started/restarted."