# auto-temp-ctrl/DEBIAN/prerm
#!/bin/bash
set -e
systemctl stop auto-temp-ctrl.service || true
systemctl disable auto-temp-ctrl.service || true
