# auto-temp-ctrl/lib/systemd/system/auto-temp-ctrl.service
[Unit]
Description=Auto Temp Control Fan Driver
After=network.target

[Service]
ExecStart=/usr/bin/fanctrl.sh
Restart=always

[Install]
WantedBy=multi-user.target
